# root/etc/s6-overlay/s6-rc.d/mt5/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid root /Metatrader/start.sh